# scripts/atlas-core.service
[Unit]
Description=Atlas Core — local orchestrator (laptop)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ronin/proyectos/atlas-core
EnvironmentFile=/home/ronin/proyectos/atlas-core/.env
ExecStart=/home/ronin/proyectos/atlas-core/.venv/bin/atlas serve --poll-interval 1.0
Restart=always
RestartSec=10
StandardOutput=append:/home/ronin/proyectos/atlas-core/.atlas.log
StandardError=append:/home/ronin/proyectos/atlas-core/.atlas.err
# Heredar variables del shell del usuario al arrancar
Environment="PYTHONPATH=/home/ronin/proyectos/atlas-core/src"

[Install]
WantedBy=default.target
